# debian/xgong.service
[Unit]
Description=Xgong server

[Service]
User=asterisk
ExecStart=/usr/bin/xgong

[Install]
WantedBy=multi-user.target
