# xgong-jenkins/debian/xgong-jenkins.cron.d
PATH=/usr/local/bin:/usr/local/sbin:/bin:/sbin:/usr/bin:/usr/sbin

* * * * 1-5 root   /usr/bin/xgong_jenkins

